# legacy/host_setup/launchd/com.gi.locate-anything.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gi.locate-anything</string>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PROJECT_DIR__/host_setup/scripts/run_locate_anything.sh</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key><false/>
        <key>Crashed</key><true/>
    </dict>
    <!-- Larger throttle than the worker (60s vs 30s) because model load
         is heavy; rapid restart loops would thrash VRAM. -->
    <key>ThrottleInterval</key><integer>60</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/gi-locate-anything.log</string>
    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/gi-locate-anything.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PROJECT_DIR__/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TZ</key>
        <string>Asia/Riyadh</string>
        <key>GI_LOCATE_WEIGHTS_DIR</key>
        <string>__USER_HOME__/Library/Caches/gi_locate</string>
    </dict>
</dict>
</plist>
