# legacy/host_setup/launchd/com.gi.cloudflared.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gi.cloudflared</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--config</string>
        <string>__USER_HOME__/.cloudflared/config.yml</string>
        <string>run</string>
        <string>gi-hub</string>
    </array>

    <key>RunAtLoad</key><true/>
    <key>KeepAlive</key><true/>
    <key>ThrottleInterval</key><integer>15</integer>

    <key>StandardOutPath</key>
    <string>__USER_HOME__/Library/Logs/gi-cloudflared.log</string>
    <key>StandardErrorPath</key>
    <string>__USER_HOME__/Library/Logs/gi-cloudflared.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TZ</key>
        <string>Asia/Riyadh</string>
    </dict>
</dict>
</plist>
